# rootfs/etc/services.d/etcdkeeper/run
#!/usr/bin/with-contenv bashio

/var/local/etcdkeeper/etcdkeeper -p 8099